# systemd/duckycap-varlink.service
[Unit]
Description=Ducky Relay Varlink Service
Documentation=man:varlinkctl(1)
Requires=duckycap-varlink.socket

[Service]
Type=notify
ExecStart=/usr/bin/duckycap-varlink --config /mnt/zed/work/youtube/scripts/ducky-relay.toml
Restart=on-failure
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ReadWritePaths=/run/duckycap.varlink

[Install]
WantedBy=multi-user.target
