# systemd/duckycap.service
[Unit]
Description=DuckyPad Capture Daemon
Documentation=file:///usr/share/doc/duckycap/README.md
# This service is activated by udev when the duckyPad is connected
# See systemd/99-duckypad.rules
Requires=duckycap-varlink.socket
After=duckycap-varlink.socket

[Service]
Type=simple
ExecStart=/usr/bin/duckycap
# Restart on failure (e.g., if device is disconnected)
Restart=on-failure
RestartSec=1
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

# Allow access to input devices
DeviceAllow=/dev/input/duckypad rw
DeviceAllow=char-input rw

[Install]
WantedBy=multi-user.target
